# etc/service/olsrd/run
#!/bin/bash -e

if [ -z "${ROUTER_MAIN_IP}" ]; then
  echo "ERROR: Router main IP not specified."
  exit 1
fi

if [ -z "${ROUTER_INTERFACES}" ]; then
  echo "ERROR: No interfaces specified."
  exit 1
fi

OLSRD_TABLE_ID=${OLSRD_TABLE_ID:-20}
OLSRD_TABLE_PRIO=${OLSRD_TABLE_PRIO:-1000}

# Configure routing policy
ip rule del from all lookup ${OLSRD_TABLE_ID} || true 2>/dev/null
ip rule add from all lookup ${OLSRD_TABLE_ID} prio ${OLSRD_TABLE_PRIO}

# Configure olsrd
OLSRD_CONFIG=/etc/olsrd/olsrd.conf
echo "AllowNoInt yes" > ${OLSRD_CONFIG}
echo "SrcIpRoutes yes" >> ${OLSRD_CONFIG}
echo "MainIp ${ROUTER_MAIN_IP}" >> ${OLSRD_CONFIG}
echo "RtTable ${OLSRD_TABLE_ID}" >> ${OLSRD_CONFIG}
echo "RtTableDefault ${OLSRD_TABLE_ID}" >> ${OLSRD_CONFIG}

echo "" >> ${OLSRD_CONFIG}
echo "LoadPlugin \"olsrd_txtinfo.so.0.1\"" >> ${OLSRD_CONFIG}
echo "{" >> ${OLSRD_CONFIG}
echo "  PlParam \"accept\" \"127.0.0.1\"" >> ${OLSRD_CONFIG}
echo "}" >> ${OLSRD_CONFIG}
echo "" >> ${OLSRD_CONFIG}

ifaces=""
for iface in ${ROUTER_INTERFACES}; do
  ifaces="${ifaces} \"${iface}\""
done
echo "Interface ${ifaces}" >> ${OLSRD_CONFIG}
echo "{" >> ${OLSRD_CONFIG}
echo "  IPv4Multicast 255.255.255.255" >> ${OLSRD_CONFIG}
echo "}" >> ${OLSRD_CONFIG}

exec /usr/sbin/olsrd -f /etc/olsrd/olsrd.conf -d 0 -nofork 2>&1

